# rootfs-k8s/files/etc/systemd/system/kubelet-wait.path
[Unit]
Description=Start kubelet when config.yaml appears

[Path]
PathExists=/var/lib/kubelet/config.yaml
Unit=kubelet.service

[Install]
WantedBy=multi-user.target
