# rootfs-k8s/files/etc/systemd/system/preload-k8s-images.service
[Unit]
Description=Preload Kubernetes container images
After=containerd.service
Before=kubelet.service
ConditionPathExists=/opt/preload-images/k8s-images.tar
ConditionPathExists=!/opt/preload-images/.loaded

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c 'ctr -n k8s.io image import /opt/preload-images/k8s-images.tar && touch /opt/preload-images/.loaded'
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
