# rootfs-k8s/files/etc/systemd/system/kubelet.service.d/00-condition.conf
[Unit]
ConditionPathExists=/var/lib/kubelet/kubeadm-flags.env
